# source/scripts/survival/openclaw-doctor.service
[Unit]
Description=Jarvis AI Physician (Emergency Recovery)

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/scripts/survival/ai_doctor.py
